# systemd/hyper-sync.service
[Unit]
Description=hyper-sync LED daemon
After=pipewire.service xdg-desktop-portal.service
PartOf=graphical-session.target

[Service]
Type=simple
Environment=__NV_PRIME_RENDER_OFFLOAD=1
Environment=__GLX_VENDOR_LIBRARY_NAME=nvidia
ExecStart=%h/.local/bin/hyper-sync daemon
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
